# etc/service/hdhomerun/run
#!/bin/bash
#

/root/hdhomerun_record_linux foreground --conf /recordings/hdhomerun.conf

